# rpi-fan.service
[Unit]
Description=Fan Controller

[Service]
Type=simple
ExecStart=/usr/local/bin/rpi-fan

[Install]
WantedBy=multi-user.target
